# upstart/flower-pi.conf
description "Service to manage flower growth"
author "Yura"

start on runlevel [2]
stop on runlevel [06]

respawn

exec flower-pi
